# docker/crontab
# LogMaster - Zamanlanmış Görevler (Docker)

# Her gün saat 02:00'da arşivleme
0 2 * * * /opt/logmaster/venv/bin/python /opt/logmaster/scripts/archiver.py >> /var/log/logmaster/archiver-cron.log 2>&1

# Her saatte bir dijital imzalama
0 * * * * /opt/logmaster/venv/bin/python /opt/logmaster/scripts/digital_signer.py --batch >> /var/log/logmaster/signer-cron.log 2>&1

# Her 6 saatte bir bütünlük kontrolü
0 */6 * * * /opt/logmaster/venv/bin/python /opt/logmaster/scripts/integrity_check.py >> /var/log/logmaster/integrity-cron.log 2>&1

# Her hafta sistem durumu raporu
0 8 * * 1 /opt/logmaster/venv/bin/python /opt/logmaster/scripts/health_check.py --report >> /var/log/logmaster/health-cron.log 2>&1

# Her ay uyumluluk raporu
0 9 1 * * /opt/logmaster/venv/bin/python /opt/logmaster/scripts/compliance_report.py --type monthly >> /var/log/logmaster/compliance-cron.log 2>&1

# Her gün log dosyalarını temizle (geçici dosyalar)
30 3 * * * find /opt/logmaster/temp -type f -mtime +1 -delete 2>/dev/null

# Her hafta eski log dosyalarını kontrol et
0 4 * * 0 find /var/log/logmaster -name "*.log" -mtime +30 -exec gzip {} \; 2>/dev/null 